# init/supervisor/serve.conf
[program:serve]
command=/project_path/goer serve
autostart=true
autorestart=true
directory=/project_path
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/project_path/storage/logs/serve.log
